# plex_debrid/services.d/webui/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
exec ttyd -p 8099 screen -x -R -S ttyd /configure.sh
